# host/llama-server/launchd/com.cyberia.llama-server-qwen36.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cyberia.llama-server-qwen36</string>

    <!--
      llama.cpp server variant for the qwen3.6-27B model A/B eval.
      Used only by host/test-model-ab/run.sh — never install this manually.

      run.sh substitutes QWEN36_GGUF_PATH at runtime via sed, then copies
      the rendered plist to ~/Library/LaunchAgents/com.cyberia.llama-server.plist
      (overwriting the production claw plist for the duration of Phase B).
      The cleanup trap restores the production plist when the run exits.

      All other flags are identical to the production claw plist so the two
      models are tested under the same conditions (same grammar, same sampler,
      same context size, same port).
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/nigel/.local/bin/llama-server</string>

        <!-- network -->
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>11435</string>

        <!-- model — substituted by run.sh from QWEN36_GGUF env var -->
        <string>--model</string>
        <string>QWEN36_GGUF_PATH</string>
        <string>--alias</string>
        <string>claw</string>

        <!-- context + memory -->
        <string>--ctx-size</string>
        <string>131072</string>
        <string>-ngl</string>
        <string>999</string>
        <string>-fa</string>
        <string>on</string>
        <string>--cache-type-k</string>
        <string>q8_0</string>
        <string>--cache-type-v</string>
        <string>q8_0</string>

        <!-- sampler — mirrors host/ollama/Modelfiles/claw.Modelfile so both
             models are tested with the same sampling conditions -->
        <string>--temp</string>
        <string>0.4</string>
        <string>--top-p</string>
        <string>0.8</string>
        <string>--top-k</string>
        <string>20</string>
        <string>--repeat-penalty</string>
        <string>1.2</string>
        <string>--repeat-last-n</string>
        <string>256</string>

        <!-- grammar — same constraint as production claw -->
        <string>--grammar-file</string>
        <string>/Users/nigel/Desktop/bench/local-llm-lab/host/llama-server/grammars/claw.gbnf</string>

        <!-- misc -->
        <string>--metrics</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/llama-server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/llama-server.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
