# host/llama-server/launchd/com.mac-llm-lab.llama-server-qwen14.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mac-llm-lab.llama-server-qwen14</string>

    <!--
      llama.cpp server variant for the Qwen3-14B Q4_K_M model A/B eval.
      Used only by host/test/run-model-qwen14-ab.sh — never install this manually.

      run-model-qwen14-ab.sh substitutes QWEN14_GGUF_PATH at runtime via sed, then copies
      the rendered plist to ~/Library/LaunchAgents/com.mac-llm-lab.llama-server.plist
      (overwriting the production claw plist for the duration of Phase B).
      The cleanup trap restores the production plist when the run exits.

      All other flags are identical to the production claw plist so the two
      models are tested under the same conditions (same grammar, same sampler,
      same context size, same port).
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/nigel/.local/bin/llama-server</string>

        <!-- network -->
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>11435</string>

        <!-- model — substituted by run-model-qwen14-ab.sh from QWEN14_GGUF env var -->
        <string>--model</string>
        <string>QWEN14_GGUF_PATH</string>
        <string>--alias</string>
        <string>claw</string>

        <!-- context + memory -->
        <string>--ctx-size</string>
        <string>32768</string>
        <string>-ngl</string>
        <string>999</string>
        <string>-fa</string>
        <string>on</string>
        <string>--cache-type-k</string>
        <string>q8_0</string>
        <string>--cache-type-v</string>
        <string>q8_0</string>
        <string>--batch-size</string>
        <string>4096</string>

        <!-- sampler — mirrors host/ollama/Modelfiles/claw.Modelfile so both
             models are tested with the same sampling conditions -->
        <string>--temp</string>
        <string>0.4</string>
        <string>--top-p</string>
        <string>0.8</string>
        <string>--top-k</string>
        <string>20</string>
        <string>--repeat-penalty</string>
        <string>1.2</string>
        <string>--repeat-last-n</string>
        <string>256</string>

        <!-- grammar — same constraint as production claw -->
        <string>--grammar-file</string>
        <string>/Users/nigel/Desktop/bench/lab/host/llama-server/grammars/claw.gbnf</string>

        <!-- misc -->
        <string>--metrics</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/llama-server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/llama-server.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
