# host/llama-server/launchd/com.mac-llm-lab.llama-server-qwen35.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mac-llm-lab.llama-server</string>

    <!--
      llama-server variant for Sprint 1.19 Exploratory Mini Sweep 2:
      Qwen3.5-9B as TIER_32 candidate.

      Substituted at install time by host/llama-server/scripts/install-qwen35
      from env vars: GGUF, CTX_SIZE, BATCH_SIZE, REPEAT_PENALTY, GRAMMAR.
      Other sampler values are frozen at the c6 (Sprint 1.19 winner) recipe
      so cells differ only on the table-axis variables.

      enable_thinking is set explicitly via --chat-template-kwargs because
      Qwen3.5 vendor docs conflict on the default and we want behavior
      comparable to the thinking-on 14B baseline.

      DELETE THIS PLIST AND THE INSTALL HELPER AFTER SWEEP 2 COALESCES.
    -->

    <key>ProgramArguments</key>
    <array>
        <string>__BINARY__</string>

        <!-- network -->
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>11435</string>

        <!-- model -->
        <string>--model</string>
        <string>__GGUF_PATH__</string>
        <string>--alias</string>
        <string>claw</string>

        <!-- context + memory -->
        <string>--ctx-size</string>
        <string>__CTX_SIZE__</string>
        <string>-ngl</string>
        <string>999</string>
        <string>-fa</string>
        <string>on</string>
        <string>--cache-type-k</string>
        <string>q8_0</string>
        <string>--cache-type-v</string>
        <string>q8_0</string>
        <string>--batch-size</string>
        <string>__BATCH_SIZE__</string>

        <!-- sampler — frozen c6 values except those varied per cell:
             - __REPEAT_PENALTY__ : Sweep 2 cells A/B/C use 1.1, D/B2 use 1.0
             - __TEMP__           : Sweep 4 cell β uses 0.7, all others 0.6
             - __MIN_P__          : Sweep 4 cell α uses 0.0 (vendor-rec),
                                    all others 0.05 (llama.cpp default).
                                    Was previously unset (implicit 0.05) for
                                    Sweep 1/2/3 apples-to-apples; now wired
                                    explicitly so cell α can flip it. -->
        <string>--temp</string>
        <string>__TEMP__</string>
        <string>--top-p</string>
        <string>0.95</string>
        <string>--top-k</string>
        <string>20</string>
        <string>--min-p</string>
        <string>__MIN_P__</string>
        <string>--repeat-penalty</string>
        <string>__REPEAT_PENALTY__</string>
        <string>--presence-penalty</string>
        <string>0.0</string>
        <string>--repeat-last-n</string>
        <string>256</string>

        <!-- thinking-mode forced on (vendor default ambiguous on Qwen3.5) -->
        <string>--chat-template-kwargs</string>
        <string>{"enable_thinking":true}</string>

        <!-- grammar — same constraint as production claw -->
        <string>--grammar-file</string>
        <string>__GRAMMAR_PATH__</string>

        <!-- misc -->
        <string>--metrics</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/llama-server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/llama-server.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
