# host/llama-server/launchd/com.home-llm-lab.llama-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.home-llm-lab.llama-server</string>

    <!--
      Template — do not install directly. Placeholders are filled by:
        ./host/llama-server/scripts/install [--size 16|32|64]

      Placeholders: __BINARY__  __GGUF_PATH__  __CTX_SIZE__  __BATCH_SIZE__  __GRAMMAR_PATH__
    -->

    <key>ProgramArguments</key>
    <array>
        <string>__BINARY__</string>

        <!-- network -->
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>11435</string>

        <!-- model -->
        <string>--model</string>
        <string>__GGUF_PATH__</string>
        <string>--alias</string>
        <string>claw</string>

        <!-- context + memory -->
        <string>--ctx-size</string>
        <string>__CTX_SIZE__</string>
        <string>-ngl</string>
        <string>999</string>
        <string>-fa</string>
        <string>on</string>
        <string>--cache-type-k</string>
        <string>q8_0</string>
        <string>--cache-type-v</string>
        <string>q8_0</string>
        <string>--batch-size</string>
        <string>__BATCH_SIZE__</string>

        <!-- sampler — mirrors host/ollama/Modelfiles/claw.Modelfile -->
        <string>--temp</string>
        <string>0.4</string>
        <string>--top-p</string>
        <string>0.8</string>
        <string>--top-k</string>
        <string>20</string>
        <!-- repeat-penalty 1.2 was over-suppressing the `\n` token in
             structured-markdown responses (it appears 8–15× in a 250-word
             reply, all inside the 256-token window). 1.05 keeps a mild
             anti-loop nudge without smushing prose. See TODO-PROSE-SMUSH.md. -->
        <string>--repeat-penalty</string>
        <string>1.05</string>
        <string>--repeat-last-n</string>
        <string>256</string>

        <!-- grammar — the point of this whole directory -->
        <string>--grammar-file</string>
        <string>__GRAMMAR_PATH__</string>

        <!-- Discipline rules: claw-code discovers CLAUDE.md in the workspace and
             concatenates it to its system prompt. See README §9. (Upstream
             llama.cpp removed `--system-prompt-file`, so per-deploy injection
             via the server is no longer available.) -->

        <!-- misc -->
        <string>--metrics</string>
    </array>

    <!-- Auto-start at login and revive on crash. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Logs go to a known place; rotate by hand if needed. -->
    <key>StandardOutPath</key>
    <string>/tmp/llama-server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/llama-server.log</string>

    <!-- Crash backoff, in case the binary dies on bad config. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
