# host/llama-server/launchd/com.mac-llm-lab.llama-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mac-llm-lab.llama-server</string>

    <!--
      Template — do not install directly. Placeholders are filled by:
        ./host/llama-server/scripts/install [--size 16|32|64]

      Placeholders:
        __BINARY__ __GGUF_PATH__ __CTX_SIZE__ __BATCH_SIZE__ __GRAMMAR_PATH__
        __TEMP__ __TOP_P__ __TOP_K__ __REPEAT_PENALTY__ __PRESENCE_PENALTY__
    -->

    <key>ProgramArguments</key>
    <array>
        <string>__BINARY__</string>

        <!-- network -->
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>11435</string>

        <!-- model -->
        <string>--model</string>
        <string>__GGUF_PATH__</string>
        <string>--alias</string>
        <string>claw</string>

        <!-- context + memory -->
        <string>--ctx-size</string>
        <string>__CTX_SIZE__</string>
        <string>-ngl</string>
        <string>999</string>
        <string>-fa</string>
        <string>on</string>
        <string>--cache-type-k</string>
        <string>q8_0</string>
        <string>--cache-type-v</string>
        <string>q8_0</string>
        <string>--batch-size</string>
        <string>__BATCH_SIZE__</string>

        <!-- sampler — per-tier values from models.conf. Flag presence is uniform
             across tiers; tier-specific values may act as no-ops (top-k=0,
             presence-penalty=0, repeat-penalty=1.0). repeat-last-n is fixed at
             256 for all tiers and bounds the window for the penalty terms. -->
        <string>--temp</string>
        <string>__TEMP__</string>
        <string>--top-p</string>
        <string>__TOP_P__</string>
        <string>--top-k</string>
        <string>__TOP_K__</string>
        <string>--repeat-penalty</string>
        <string>__REPEAT_PENALTY__</string>
        <string>--presence-penalty</string>
        <string>__PRESENCE_PENALTY__</string>
        <string>--repeat-last-n</string>
        <string>256</string>

        <!-- grammar — the point of this whole directory -->
        <string>--grammar-file</string>
        <string>__GRAMMAR_PATH__</string>

        <!-- Discipline rules: claw-code discovers CLAUDE.md in the workspace and
             concatenates it to its system prompt. See README §9. (Upstream
             llama.cpp removed `--system-prompt-file`, so per-deploy injection
             via the server is no longer available.) -->

        <!-- ====================================================================
             FUTURE-MODEL HINT — chat-template-kwargs (Sprint 1.19, 2026-05-01)

             Forces `enable_thinking: true` for hybrid-thinking models whose
             default is ambiguous (e.g. Qwen3.5-9B at tier-16, where the model
             card and llama.cpp Jinja template disagree on the default). Without
             this flag, cell E (Sweep 2/3 tier-16 winner) silently degrades.

             Currently safe across all tiers because:
               - tier-16 (Qwen3.5-9B IQ4_XS): designed for this flag (cells E/γ/δ/ζ).
               - tier-32 (Qwen3.5-9B UD-Q5_K_XL @ 64k, B2): same family as
                 tier-16 — designed for this flag.
               - tier-64 (Qwen3.6-35B): default-thinking, **suppressed via litellm
                 per-request override** (see host/litellm/litellm-config.yaml
                 `claw-llama` extra_body.chat_template_kwargs.enable_thinking=false).
                 The per-request override wins against the server-side default.

             FUTURE: if a model lands whose Jinja template ignores
             `enable_thinking`, OR whose template uses a different kwarg name
             (e.g. `thinking`, `reasoning_mode`), this flag may become a no-op
             or, worse, cause a template-render error. Re-evaluate when:
               1. Adopting a new model family at any tier (re-test the install).
               2. Upgrading llama.cpp past a release-note that mentions
                  chat-template-kwargs handling changes.
               3. Adding a tier whose chat template doesn't support
                  enable_thinking.

             Sweep 2 cells used a separate variant plist (com.mac-llm-lab.llama-server-qwen35.plist)
             before this flag was promoted to production. That variant is now
             redundant for tier-16 in steady-state — kept on disk as a reference
             implementation and a fallback if this universal flag has to be
             reverted.
        ==================================================================== -->
        <string>--chat-template-kwargs</string>
        <string>{"enable_thinking":true}</string>

        <!-- misc -->
        <string>--metrics</string>
    </array>

    <!-- Auto-start at login and revive on crash. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Logs go to a known place; rotate by hand if needed. -->
    <key>StandardOutPath</key>
    <string>/tmp/llama-server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/llama-server.log</string>

    <!-- Crash backoff, in case the binary dies on bad config. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
