# host/llama-server/launchd/com.mac-llm-lab.llama-server-baseline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mac-llm-lab.llama-server</string>

    <!--
      BASELINE config for settings A/B test (run-model-settings-ab.sh Phase A).
      Differences from the production plist (com.mac-llm-lab.llama-server.plist):
        - --ctx-size 131072  (production: 32768)
        - no --batch-size    (production: 4096)
      Everything else is identical, including the corrected grammar-file path.
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/nigel/.local/bin/llama-server</string>

        <!-- network -->
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>11435</string>

        <!-- model -->
        <string>--model</string>
        <string>/Users/nigel/.ollama/gguf/Qwen3-Coder-30B-A3B-Instruct-UD-Q6_K_XL.gguf</string>
        <string>--alias</string>
        <string>claw</string>

        <!-- context + memory — BASELINE values -->
        <string>--ctx-size</string>
        <string>131072</string>
        <string>-ngl</string>
        <string>999</string>
        <string>-fa</string>
        <string>on</string>
        <string>--cache-type-k</string>
        <string>q8_0</string>
        <string>--cache-type-v</string>
        <string>q8_0</string>
        <!-- no --batch-size: llama-server default is 2048 -->

        <!-- sampler — mirrors host/ollama/Modelfiles/claw.Modelfile -->
        <string>--temp</string>
        <string>0.4</string>
        <string>--top-p</string>
        <string>0.8</string>
        <string>--top-k</string>
        <string>20</string>
        <string>--repeat-penalty</string>
        <string>1.2</string>
        <string>--repeat-last-n</string>
        <string>256</string>

        <!-- grammar -->
        <string>--grammar-file</string>
        <string>/Users/nigel/Desktop/bench/lab/host/llama-server/grammars/claw.gbnf</string>

        <!-- misc -->
        <string>--metrics</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/llama-server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/llama-server.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
